# config/crontab
*/10 * * * * wget -q --server-response --post-data="" "https://${DRAFT_MODE_HOSTNAME}/api/memoized-fns-invalidate?token=XXX" -O - && echo